# obsidian-autocommit.service
[Unit]
Description=Auto-commit Obsidian changes to Git
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/jbi-laptop/Git/009-baena.blog-personal-pkm
ExecStart=/home/jbi-laptop/Git/009-baena.blog-personal-pkm/auto-commit.sh
User=jbi-laptop
StandardOutput=append:/tmp/obsidian-autocommit.log
StandardError=append:/tmp/obsidian-autocommit.log
